# systemd/ai-employee-reflect.timer
# Fires ai-employee-reflect.service weekly, Sunday ~22:00 local time.
# Enable with:  sudo systemctl enable --now ai-employee-reflect.timer
[Unit]
Description=Self-hosted AI employee — weekly reflection timer

[Timer]
OnCalendar=Sun *-*-* 22:00:00
# if the box was asleep/off at the scheduled time, run as soon as it's back
Persistent=true
# small jitter so it doesn't hammer providers at exactly :00
RandomizedDelaySec=300
Unit=ai-employee-reflect.service

[Install]
WantedBy=timers.target
