# systemd/ai-employee-digest.service
# Self-hosted AI employee — daily morning brief (one-shot, fired by the .timer).
# EDIT the two paths/user below to match ai-employee.service, then:
#   sudo cp systemd/ai-employee-digest.service systemd/ai-employee-digest.timer /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now ai-employee-digest.timer
#   systemctl list-timers ai-employee-digest.timer        # check next run
#   journalctl -u ai-employee-digest -f                   # watch a run
[Unit]
Description=Self-hosted AI employee — daily digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=youruser
WorkingDirectory=/home/youruser/ai-employee
ExecStart=/home/youruser/ai-employee/.venv/bin/python -m agent digest
# treat the agent as untrusted — OS isolation is the real safety net
MemoryMax=2G
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/youruser/ai-employee
