# systemd/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway (Telegram, Gemini)
After=network-online.target ollama.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=hermes
Group=hermes
WorkingDirectory=/home/hermes
Environment=HERMES_ACCEPT_HOOKS=1
ExecStart=/home/hermes/.local/bin/hermes gateway run
Restart=on-failure
RestartSec=10
MemoryMax=2G
CPUQuota=200%
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/hermes
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
