# systemd/hermes-pipeline-weekly.timer
[Unit]
Description=weekly maintenance pipeline (nurture sweep + vacuum) - Sunday 21:30
[Timer]
OnCalendar=Sun *-*-* 21:30:00
Persistent=true
[Install]
WantedBy=timers.target
